# config/features/sourceSystemPatching/crontab
29 4 * * * /bin/ksh /var/cfengine/inputs/config/features/sourceSystemPatching/update-src
