# config/features/databaseServer/crontab
0 0 * * * svc -t /service/postgresql/log
