# config/features/samba/crontab
0	0	*	*	*	svc -t /service/smbd/log
0 	*	*	*	*	/bin/ksh /var/cfengine/inputs/config/features/samba/mail_smbd_logs.sh
