# config/features/sourcePackageAdministration/crontab
29 4 * * * /bin/ksh /var/cfengine/inputs/config/features/sourcePackageAdministration/update-ports
