# config/features/kerberosSlave/kerberos_slave.crontab
15	*	*	*	*	/bin/sh /var/cfengine/inputs/config/features/kerberosSlave/mirror_db

