# config/features/timeSynchronization/rc.local

(sleep 10 && rdate -ns ntp1.belbone.be) &

